# ikuta_am.service
[Unit]
Description=IKUTA Annouce Machine
After=network.target

[Service]
ExecStart=/opt/ikuta_am/ikuta_am.sh
WorkingDirectory=/opt/ikuta_am
