# zabbix-agent-restart.service
[Unit]
Description=Scheduled restart zabbix agent

[Service]
Type=simple
ExecStart=/bin/systemctl restart zabbix-agent
