# zabbix-agent-restart.timer
[Unit]
Description=Run zabbix-agent-restart monthly

[Timer]
OnBootSec=1M
OnUnitActiveSec=1M

[Install]
WantedBy=timers.target
